# ice-turn/build/client.service
[Unit]
After=network.target
Description=ICE TURN Client to connect with ICE TURN Connector
[Install]
WantedBy=multi-user.target
[Service]
Type=simple
ExecStart=/usr/local/bin/ice-turn-client service start --log-proxy-throughput 1 --log-throughput-endpoint http://13.215.59.103:8888/hop-results
ExecStop=/usr/local/bin/ice-turn-client a stop
WorkingDirectory=/usr/local/bin
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n% 